# dot_config/systemd/user/chezmoi-auto-save.path
[Unit]
Description=Watch pacman db

[Path]
PathModified=/var/lib/pacman/local

[Install]
WantedBy=paths.target
