# dot_config/systemd/user/chezmoi-watch-packages.path
[Unit]
Description=Watch pacman db for package changes

[Path]
PathModified=/var/lib/pacman/local

[Install]
WantedBy=paths.target
